# conf/rustdesksignal.service
[Unit]
Description=Rustdesk Signal Server

[Service]
Type=simple
User=__APP__
Group=__APP__
Restart=always
LimitNOFILE=1000000
WorkingDirectory=__INSTALL_DIR__/
ExecStart=__INSTALL_DIR__/hbbs -r 127.0.0.1 -p __PORT_HBBS__
StandardOutput=append:/var/log/__APP__/hbbs.log
StandardError=append:/var/log/__APP__/hbbs.error
RestartSec=10

[Install]
WantedBy=multi-user.target